# homework-7/watcher/watcher.service
[Unit]
Description=My watcher service

[Service]
Type=oneshot
EnvironmentFile=/etc/sysconfig/watcher.conf
ExecStart=/opt/watcher.sh $WORD $LOG
